# WorkingProject/test.delegator.service
[Unit]
Description=Delegation Service
[Service]
Type=dbus
BusName=test.delegator.StandAlone
ExecStart=/home/ubuntu/delegate.c
